# upsMetrics.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * pi /home/pi/ep2000-cloudwatch/run.sh > /home/pi/ep2000-cloudwatch/last-run.log 2>&1
